# conf/memcached.ini
[program:memcached]
command=/usr/bin/memcached -p 11211 -l 127.0.0.1 -u nobody -m 512m
autostart=true
autorestart=true
stdout_logfile=/proc/self/fd/1
stderr_logfile=/proc/self/fd/2
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
